# profile.d/grunt-completion.bash
#!/bin/bash

hash grunt 2>/dev/null && \
eval "$(grunt --completion=bash)"
