# profile.d/grunt-completion.bash
#!/bin/bash

hash grunt >/dev/null && \
eval "$(grunt --completion=bash)"
