# profile.d/podman.bash
#!/bin/bash

# hash podman 2>/dev/null && \
# alias docker="podman"

